# etc/systemd/system/getty@tty1.service.d/override.conf
# Automatic login of user markus (which is safe as the encryption password has been entered before)
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin markus --noclear %I $TERM
